# shots/ws/vmscreenshot.service
[Unit]
Description="Сделать снимок экрана работающей ВМ"

[Service]
ExecStart=/opt/scripts/vmgetshot
